# deploy/bc-hearing-watch.service
[Unit]
Description=BC Local Government Council Tracker — FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/lg-tracking/backend
ExecStart=/var/www/lg-tracking/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info
Restart=always
RestartSec=5

# Load environment variables from file
EnvironmentFile=/var/www/lg-tracking/backend/.env

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
